# surfduck/surfduck.service
[Unit]
Description=SurfDuck — curated TV guide. Find what's worth watching.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omar
Group=omar
WorkingDirectory=/home/omar/agentos/surfduck
EnvironmentFile=/home/omar/agentos/surfduck/.env
ExecStart=/usr/bin/python3 /home/omar/agentos/surfduck/server.py
Restart=always
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/omar/agentos/surfduck/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target